# distros/ubuntu/config/etc/profile.d/brew.sh
#!/bin/bash

if command -v /home/linuxbrew/.linuxbrew/bin/brew &> /dev/null; then
    eval $(/home/linuxbrew/.linuxbrew/bin/brew shellenv)
fi
